# serverwatch.service
[Unit]
Description=ServerWatch - Minecraft Server Monitoring
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/serverwatch
Environment="PATH=/usr/bin:/usr/local/bin"
ExecStart=/usr/bin/python3 /opt/serverwatch/main.py
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=serverwatch

[Install]
WantedBy=multi-user.target
